# simpleService.service
[Unit]
Description=Service ARP snapshot
After=network.target

[Service]
ExecStart=/etc/simpleService.sh